# S99git_backup
#!/bin/sh /etc/rc.common

START=99
STOP=1
PROG=/usr/local/bin/git_backup
PIDFILE=/var/run/git_backup.pid
LOGFILE=/var/log/git_backup.log

start() {
    echo "Starting git backup service..."
    # Create a wrapper that handles logging
    start-stop-daemon -S -b -m -p $PIDFILE -x /bin/sh -- -c "$PROG >> $LOGFILE 2>&1"
}

stop() {
    echo "Stopping git backup service..."
    # -K: stop (kill)
    # -q: quiet (suppress common errors)
    # -p: use this PID file to find the process
    # -s TERM: use SIGTERM for graceful shutdown first
    start-stop-daemon -K -q -p $PIDFILE -s TERM
    
    # Wait a moment for graceful shutdown
    sleep 2
    
    # Force kill if still running
    start-stop-daemon -K -q -p $PIDFILE -s KILL 2>/dev/null
    
    # Clean up PID file
    rm -f $PIDFILE
}
